# conf/tailscale.conf
[program:tailscale]
command     = /root/start.sh
user        = root

autostart=true
autorestart=true
startsecs=60
stopasgroup=true
ikillasgroup=true
startretries=3
redirect_stderr         = true
stdout_logfile_maxbytes = 50MB
stdout_logfile_backups  = 10
stdout_logfile          = /var/log/supervisor/tailscale.log